# supervisord.conf
[supervisord]
nodaemon=true
user=root

[program:node]
command=node server.js
directory=/app
autostart=true
autorestart=true