# supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisord.log
pidfile=/var/run/supervisord.pid

[program:qbittorrent]
command=qbittorrent-nox
autostart=true
autorestart=true
stderr_logfile=/var/log/qbittorrent.err.log
stdout_logfile=/var/log/qbittorrent.out.log
stopwaitsecs=20
startsecs=10
priority=100

[program:node]
command=node server.js
directory=/app
autostart=true
autorestart=true
stderr_logfile=/var/log/node.err.log
stdout_logfile=/var/log/node.out.log
stopwaitsecs=30
startsecs=20
priority=200
environment=NODE_ENV=production,PORT=3000

[eventlistener:processes]
command=sh -c "echo READY && read line && kill -SIGTERM $PPID"
events=PROCESS_STATE_STOPPED,PROCESS_STATE_EXITED,PROCESS_STATE_FATAL
autorestart=true